# e2e/templates/discovery.service.j2
[Unit]
Description=Prime Network Discovery
After=network.target

[Service]
Type=simple
User={{ ansible_user }}
Environment=VALIDATOR_ADDRESS={{ validator_address }}
Environment=RPC_URL={{ rpc_url }}
WorkingDirectory=/opt/prime/discovery
ExecStart=/opt/prime/discovery/discovery --validator-address ${VALIDATOR_ADDRESS} --rpc-url ${RPC_URL}

[Install]
WantedBy=multi-user.target