# e2e/templates/validator.service.j2
[Unit]
Description=Prime Network Validator
After=network.target

[Service]
Type=simple
User={{ ansible_user }}
Environment=VALIDATOR_KEY={{ validator_key }}
Environment=RPC_URL={{ rpc_url }}
Environment=DISCOVERY_URL={{ discovery_url }}
WorkingDirectory=/opt/prime/validator
ExecStart=/opt/prime/validator/validator --validator-key ${VALIDATOR_KEY} --rpc-url ${RPC_URL} --discovery-url ${DISCOVERY_URL}
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target